# ops/systemd/daily-lifecycle-report.service
[Unit]
Description=Generate daily lifecycle report
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/srv/automation-mvp
ExecStart=/srv/automation-mvp/.venv/bin/python daily_lifecycle_report.py --since-hours 24 --output-md /srv/automation-mvp/data/reports/daily_lifecycle_report_latest.md --output-json /srv/automation-mvp/data/reports/daily_lifecycle_report_latest.json
