# ops/systemd/edge-search-mutation-cycle.timer
[Unit]
Description=Run edge search mutation cycle every 30 minutes

[Timer]
OnCalendar=*:0/30
Persistent=true
Unit=edge-search-mutation-cycle.service

[Install]
WantedBy=timers.target
