# ops/systemd/edge-search-manifest-worker.service
[Unit]
Description=Edge search manifest worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/automation-mvp
EnvironmentFile=-/srv/automation-mvp/ops/systemd/edge-search.env
ExecStart=/srv/automation-mvp/.venv/bin/python manifest_worker.py
Restart=always
RestartSec=5
